# systemdUser/.config/systemd/user/wallpaper.timer
[Unit]
Description=run my wallpaper script every 5 minutes

[Timer]
OnBootSec=5m
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
